# src/resources/com.fortivpn-tray.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fortivpn-tray.helper</string>
    <key>Program</key>
    <string>/Library/PrivilegedHelperTools/fortivpn-helper</string>
    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SockPathName</key>
            <string>/var/run/fortivpn-helper.sock</string>
            <key>SockMode</key>
            <integer>438</integer>
        </dict>
    </dict>
    <key>KeepAlive</key>
    <false/>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.ktutnik.fortivpn-tray</string>
    </array>
</dict>
</plist>
